# .root/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --nonewline --noissue --skip-login --autologin khuedoan --noclear %I $TERM
TTYVTDisallocate=no
